# PhotoBooth.conf
[program:photo-booth]
command=python /home/pi/PhotoBooth/PhotoBooth.py
autostart=true
autorestart=true
user=root